# profile.d/6-editor.sh
# editor.sh
export EDITOR=/usr/bin/nvim
